# paulshaclaw/memory/dream/systemd/paulsha-memory-dream.timer
[Unit]
Description=Run PaulSha memory dream on workday mornings

[Timer]
OnCalendar=Mon..Fri 05:00
Persistent=true

[Install]
WantedBy=timers.target
