# .test-artifacts/stage7-install-create-only/home/.config/systemd/user/demo-agent-cost.service
[Unit]
Description=demo-agent cost refresh service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
EnvironmentFile=%h/.agents/core/runtime/demo-agent.env
EnvironmentFile=%h/.agents/core/runtime/demo-agent-cost.env
ExecStart=/usr/bin/env bash /srv/paulshaclaw/scripts/service-cost.sh
KillMode=control-group
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
