# .test-artifacts/stage7-install-create-only/home/.config/systemd/user/demo-agent-manager.service
[Unit]
Description=demo-agent persona manager service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
EnvironmentFile=%h/.agents/core/runtime/demo-agent.env
EnvironmentFile=%h/.agents/core/runtime/demo-agent-manager.env
ExecStart=/usr/bin/env bash /srv/paulshaclaw/scripts/service-manager.sh
KillMode=control-group
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
