# paulshaclaw/deploy/templates/core/systemd/__INSTANCE__.service.tmpl
[Unit]
Description=__INSTANCE__ service

[Service]
ExecStart=/usr/bin/env python3 -m paulshaclaw.core.daemon

[Install]
WantedBy=multi-user.target
