# paulshaclaw/deploy/templates/core/systemd/__INSTANCE__-telegram.service.tmpl
[Unit]
Description=__INSTANCE__ Telegram listener
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=%h/.agents/core/runtime/__INSTANCE__.env
EnvironmentFile=%h/.agents/core/runtime/__INSTANCE__-telegram.env
EnvironmentFile=%h/.config/paulshaclaw/__INSTANCE__.telegram.secret.env
Environment=PSC_STAGE1_CONFIG=%h/.agents/state/config/__INSTANCE__.state.json
ExecStart=/usr/bin/env bash __ROOT_DIR__/scripts/service-bot.sh
KillMode=control-group
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target
