# paulshaclaw/deploy/templates/core/systemd/__INSTANCE__-manager.service.tmpl
[Unit]
Description=__INSTANCE__ persona manager tick (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PSC_MANAGER_SPECS_DIR=%h/.agents/specs
EnvironmentFile=%h/.agents/core/runtime/__INSTANCE__.env
EnvironmentFile=%h/.agents/core/runtime/__INSTANCE__-manager.env
Environment=PSC_STAGE1_CONFIG=%h/.agents/state/config/__INSTANCE__.state.json
ExecStart=/usr/bin/env python3 -m paulshaclaw.coordinator tick --require-idle --specs-dir ${PSC_MANAGER_SPECS_DIR} --executor ${PSC_MANAGER_EXECUTOR}
