# paulshaclaw/deploy/templates/core/systemd/__INSTANCE__-manager.timer.tmpl
[Unit]
# DEPRECATED: manager 現改為常駐 service；保留模板僅供舊部署辨識/回滾。
Description=__INSTANCE__ persona manager tick timer (DEPRECATED)

[Timer]
OnBootSec=2min
OnUnitActiveSec=300
Unit=__INSTANCE__-manager.service

[Install]
WantedBy=timers.target
