# paulshaclaw/deploy/templates/core/systemd/__INSTANCE__-manager.service.tmpl
[Unit]
Description=__INSTANCE__ persona manager service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
EnvironmentFile=%h/.agents/core/runtime/__INSTANCE__.env
EnvironmentFile=%h/.agents/core/runtime/__INSTANCE__-manager.env
ExecStart=/usr/bin/env bash __ROOT_DIR__/scripts/service-manager.sh
KillMode=control-group
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
