# paulshaclaw/deploy/templates/core/systemd/__INSTANCE__-manager.timer.tmpl
# Deprecated: manager now runs as a persistent service; keep this timer only for catalog/rollback reference.
[Unit]
Description=Deprecated __INSTANCE__ persona manager timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=300
Unit=__INSTANCE__-manager.service

[Install]
WantedBy=timers.target
