# paulshaclaw/deploy/templates/core/systemd/__INSTANCE__-manager.timer.tmpl
[Unit]
Description=__INSTANCE__ persona manager tick timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=300
Unit=__INSTANCE__-manager.service

[Install]
WantedBy=timers.target
